# mikrokredit.service
[Unit]
Description=MikroKredit Web Application
After=network.target postgresql.service

[Service]
Type=simple
User=valstan
WorkingDirectory=/home/valstan/development/mikrokredit
Environment=MIKROKREDIT_DATABASE_URL=postgresql://mikrokredit_user:mikrokredit_pass_2024@localhost:5434/mikrokredit_db
ExecStart=/home/valstan/development/mikrokredit/.venv/bin/python -c "from web import create_app; app = create_app(); app.run(host='0.0.0.0', port=8002, debug=False)"
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target